# dist/systemd/linux-status.service
[Unit]
Description=linux-status

[Service]
Type=simple
WorkingDirectory=/srv/linux-status
ExecStart=/srv/linux-status/server.py

[Install]
WantedBy=multi-user.target
