# dist/systemd/linux-status.service
[Unit]
Description=linux-status

[Service]
Type=simple
WorkingDirectory=/srv/linux-status
ExecStart=/usr/bin/python3 -B server.py

[Install]
WantedBy=multi-user.target
